# efkstack/config/elasticsearch-supervisor.ini
[program:elasticsearch]
user=elasticsearch
command=/build/elasticsearch/bin/elasticsearch
