# efkstack/config/fluentd-supervisor.ini
[program:fluentd]
command=/usr/sbin/td-agent
